# crontab
0 0 * * * /app/venv/bin/python /app/failman.py >> /proc/1/fd/1 2>> /proc/1/fd/2
